# docker/crontab
30 * * * * source /etc/profile; cd /var/www/html/ && /usr/local/bin/php index.php -j=refresh
30 1 * * * source /etc/profile; cd /var/www/html/ && /usr/local/bin/php index.php -j=clear-log
